# agents/shiftshapr/com.agentlab.shiftshapr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agentlab.shiftshapr</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/shiftshapr/workspace/agent-lab/.venv/bin/python</string>
    <string>/Users/shiftshapr/workspace/agent-lab/agents/shiftshapr/shiftshapr.py</string>
    <string>--port</string>
    <string>8080</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/shiftshapr/workspace/agent-lab</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/shiftshapr/workspace/agent-lab/logs/shiftshapr.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/shiftshapr/workspace/agent-lab/logs/shiftshapr.err</string>
</dict>
</plist>
